# LaunchAgent/com.maxweisel.maxvoice.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.maxweisel.maxvoice</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>
MARKER="$HOME/.maxvoice/.accessibility_failed"
if [ -f "$MARKER" ]; then
    echo "Accessibility permission not granted, not starting MaxVoice"
    exit 1
fi
exec /Applications/MaxVoice.app/Contents/MacOS/MaxVoice
        </string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardErrorPath</key>
    <string>/tmp/maxvoice.err</string>

    <key>StandardOutPath</key>
    <string>/tmp/maxvoice.log</string>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
